# supervisor/conf.d/webmin.conf
[program:webmin]
priority=1
command=/etc/webmin/start
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
autorestart=false
startsecs=0
